# ops/com.nickltang.sf-apartment-aggregator.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nickltang.sf-apartment-aggregator.poll</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__REPO_DIR__/scripts/run_poll.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>StartInterval</key>
  <integer>1200</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/logs/poll.out.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/logs/poll.err.log</string>
</dict>
</plist>
